# aemet-exporter.timer
[Unit]
Description=AEMET Exporter Timer

[Timer]
OnCalendar=*-*-* *:20:00

[Install]
WantedBy=timers.target