# aemet-exporter.timer
[Unit]
Description=AEMET Exporter Timer

[Timer]
OnCalendar=Hourly

[Install]
WantedBy=timers.target